# imgvidgen.service
# systemd 유닛 — sudo cp 로 /etc/systemd/system/imgvidgen.service 설치 후 enable --now.
# imgvidgen.ai-ve.uk (Cloudflare hyfe-iqc 터널) -> 127.0.0.1:8503 로 프록시됨.
[Unit]
Description=Image & Video Generator (imgvidgen.ai-ve.uk)
After=network-online.target
Wants=network-online.target
# 재시작 횟수 제한 해제: 빠른 크래시-루프(기본 10초 내 5회)에 systemd 가
# 영구 포기하면 502 가 고착된다. 죽으면 무조건 다시 띄운다.
StartLimitIntervalSec=0

[Service]
Type=simple
User=opc
Group=opc
WorkingDirectory=/home/opc/projects/image_generator
Environment=APP_PORT=8503
# 서브도메인 루트(/)에서 서빙되므로 Gradio root_path 를 비워 /p/<port> 프리픽스를 끈다.
Environment=GRADIO_ROOT_PATH=
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
